# srv/images/cron-provider1/etc/crontabs/root
# note: run `npm run db-setup` first to initalize tables

# HOME=~/projects/air_quality/IST VM/srv/images/cron-provider1/

# crontab runs from /root by default
# use this first: `cd /srv`


# sync purpleair sensor data
*/30 * * * * cd /srv && /usr/bin/npm run pa-sync > /proc/1/fd/1 2>/proc/1/fd/2

# sync FEM station data
20,30,40 * * * * cd /srv && /usr/bin/npm run station-sync > /proc/1/fd/1 2>/proc/1/fd/2


# TODO: sync metadata